# roles/video_downloader/files/ytdl-sub.timer
[Unit]
Description=Download YouTube Subscriptions

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=12

[Install]
WantedBy=timers.target
